# conf/etc/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=uwsgi --master /etc/uwsgi/whiteboard.ini
directory=/deploy/code
autostart=false